# wifi-auto-catch-handshakes/wifi-auto-catch-handshakes.service
[Unit]
Description=pentest: catch wifi handshakes

[Service]
ExecStart=/usr/bin/npm start /home/pi/wifi-auto-catch-handshakes/
User=root
Restart=always
WorkingDirectory=/home/pi/wifi-auto-catch-handshakes

[Install]
WantedBy=multi-user.target