# wifi-auto-capture-open-traffic/wifi-auto-capture-open-traffic.service
[Unit]
Description=capture open wifi traffic

[Service]
ExecStart=/usr/bin/npm start /home/pi/wifi-auto-capture-open-traffic/
User=root
Restart=always
WorkingDirectory=/home/pi/wifi-auto-capture-open-traffic

[Install]
WantedBy=multi-user.target